# templates/default/uptime-status.conf.erb
description "Uptime Status App"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

chdir <%= @dir %>
setuid <%= @user %>

script
  export NODE_ENV=status
  exec node status.js >> <%= @log_dir %>/uptime-status.log 2>&1
end script

emits uptime-status-running
